# supervisord.conf
[supervisord]
http_port=/var/tmp/supervisor.sock ; (default is to run a UNIX domain socket server)
user=root
logfile=/var/log/supervisor/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB       ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10          ; (num of main logfile rotation backups;default 10)
loglevel=info               ; (logging level;default info; others: debug,warn)
pidfile=/var/run/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false              ; (start in foreground if true;default false)
minfds=1024                 ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)

[inet_http_server]         ; inet (TCP) server disabled by default
port=*:9001        ; (ip_address:port specifier, *:port for all iface)
;username=user              ; (default is no username (open server))
;password=123               ; (default is no password (open server))

[unix_http_server]
file=/var/tmp/supervisor.sock   ; (the path to the socket file)
chmod=0770

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/tmp/supervisor.sock ; use a unix:// URL  for a unix socket

; The below sample program section shows all possible program subsection values,
; create one or more 'real' program: sections to be able to control them under
; supervisor.

[program:mailer]
command=/usr/local/bin/python3.4 /root/techaspects/backend/tornado/mailer.py
user=root
autorestart=true
autostart=true
startsecs=10

[program:readdressweb]
command=/usr/local/bin/python3.4 /root/mainrepo/app.py
user=root
autorestart=true
autostart=true
startsecs=10

[program:readdressehndlr]
command=/usr/local/bin/python3.4 /root/mainrepo/emailhandler.py -i %(process_num)02d
user=root
autorestart=true
autostart=true
startsecs=10
numprocs=4
process_name=%(program_name)s_%(process_num)s

[program:readdressahndlr]
command=/usr/local/bin/python3.4 /root/mainrepo/archivemail.py -i %(process_num)02d
user=root
autorestart=true
autostart=true
startsecs=10
numprocs=4
process_name=%(program_name)s_%(process_num)s

[program:readdressshndlr]
command=/usr/local/bin/python3.4 /root/mainrepo/sendmail.py -i %(process_num)02d
user=root
autorestart=true
autostart=true
startsecs=10
numprocs=4
process_name=%(program_name)s_%(process_num)s

[program:readdressshndlr]
command=/usr/local/bin/python3.4 /root/mainrepo/sendmail.py -i %(process_num)02d
user=root
autorestart=true
autostart=true
startsecs=10
numprocs=4
process_name=%(program_name)s_%(process_num)s

[program:readdressdhndlr]
command=/usr/local/bin/python3.4 /root/mainrepo/deregistermailid.py -i %(process_num)02d
user=root
autorestart=true
autostart=true
startsecs=10
numprocs=4
process_name=%(program_name)s_%(process_num)s

[program:readdressghndlr]
command=/usr/local/bin/python3.4 /root/mainrepo/genbouncemail.py -i %(process_num)02d
user=root
autorestart=true
autostart=true
startsecs=10
numprocs=4
process_name=%(program_name)s_%(process_num)s

[program:readdresslhndlr]
command=/usr/local/bin/python3.4 /root/mainrepo/liarchive.py -i %(process_num)02d
user=root
autorestart=true
autostart=true
startsecs=10
numprocs=4
process_name=%(program_name)s_%(process_num)s

[program:readdressmhndlr]
command=/usr/local/bin/python3.4 /root/mainrepo/modifymailid.py -i %(process_num)02d
user=root
autorestart=true
autostart=true
startsecs=10
numprocs=4
process_name=%(program_name)s_%(process_num)s

